# conf/crontab
*/5 * * * * /usr/local/bin/proxyctl cron
