# conf/crontab
*/15 * * * * /usr/local/bin/stop_projects.sh
