# conf/crontab
# Run `proxyctl cron` every minute
* * * * * /usr/local/bin/proxyctl cron

# Run `proxyctl cleanup` every hour
0 * * * * /usr/local/bin/proxyctl cleanup
